# salt/opentsdb/templates/opentsdb.service.tpl
[Unit]
Description=opentsdb

[Service]
Type=simple
ExecStart=/usr/bin/tsdb tsd